# etc/udev/rules.d/91-usb-swap.rules
SUBSYSTEMS=="usb", RUN+="/opt/kvm/kvm.sh"